# apps/frontend/scripts/auto-healing/systemd/n8n-auto-healing-monitor.service
[Unit]
Description=n8n Auto-Healing Monitor Service
Documentation=https://github.com/YOUR_REPO/gonsai2
After=network-online.target docker.service
Wants=network-online.target
Requires=n8n-auto-healing-monitor.timer

[Service]
Type=oneshot
User=%i
Group=%i
WorkingDirectory=/home/%i/projects/gonsai2/apps/frontend

# 환경 변수
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/%i/.nvm/versions/node/v20.11.0/bin"
Environment="NODE_ENV=production"
EnvironmentFile=-/home/%i/projects/gonsai2/apps/frontend/.env.local

# 스크립트 실행
ExecStart=/bin/bash /home/%i/projects/gonsai2/apps/frontend/scripts/auto-healing/monitor.sh

# 로깅
StandardOutput=journal
StandardError=journal
SyslogIdentifier=n8n-auto-healing-monitor

# 보안 설정
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/projects/gonsai2/apps/frontend/scripts/auto-healing/logs
ReadWritePaths=/home/%i/projects/gonsai2/apps/frontend/scripts/auto-healing/state
ReadWritePaths=/home/%i/projects/gonsai2/apps/frontend/scripts/auto-healing/backups

# 리소스 제한
CPUQuota=50%
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
