# apps/frontend/scripts/auto-healing/systemd/n8n-auto-healing-monitor.timer
[Unit]
Description=n8n Auto-Healing Monitor Timer
Documentation=https://github.com/YOUR_REPO/gonsai2
Requires=n8n-auto-healing-monitor.service

[Timer]
# 5분마다 실행
OnBootSec=5min
OnUnitActiveSec=5min

# 정확한 시간에 실행하지 않고 약간의 랜덤성 추가 (부하 분산)
RandomizedDelaySec=30sec

# 시스템이 절전 모드에서 깨어났을 때 놓친 작업 실행
Persistent=true

[Install]
WantedBy=timers.target
